# config/com.blue-pearmain.pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.blue-pearmain.pipeline.plist
  launchd agent for Blue Pearmain full maintenance run

  Runs 'bp all' every 6 hours:
    1. scan --all       full Apple Photos index refresh
    2. poll             catch recent Flickr uploads
    3. thumbs           populate missing thumbnail paths
    4. pipeline         metadata sync + auto-apply non-conflict proposals
    5. reconcile --fix  repair DB/Flickr mismatches
    6. sync-albums      push album memberships to Flickr
    7. checkpoint       trim the WAL file

  Each step runs independently; a failure is logged and the sequence continues.
  Collision proposals are left pending for manual review in the UI.

  Install with bp (substitutes __REPO__, __UV__, __HOME__ automatically):
    bp install-daemons

  Or install manually after substituting tokens:
    launchctl bootstrap gui/$(id -u) \
       ~/Library/LaunchAgents/com.blue-pearmain.pipeline.plist

  Run now (without waiting for interval):
    launchctl start com.blue-pearmain.pipeline

  Uninstall:
    launchctl bootout gui/$(id -u) \
       ~/Library/LaunchAgents/com.blue-pearmain.pipeline.plist
    rm ~/Library/LaunchAgents/com.blue-pearmain.pipeline.plist

  Logs:
    tail -f __HOME__/Library/Logs/BluePearmain/pipeline.log

  NOTE: Photos.app must be running for proposals targeting Apple Photos to apply.
  Proposals are left pending if Photos is not open; they will apply on the next run.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.blue-pearmain.pipeline</string>

    <key>ProgramArguments</key>
    <array>
        <string>__UV__</string>
        <string>run</string>
        <string>python</string>
        <string>__REPO__/bp</string>
        <string>all</string>
        <string>--config</string>
        <string>__REPO__/config/config.yml</string>
    </array>

    <!-- Run every 6 hours -->
    <key>StartInterval</key>
    <integer>21600</integer>

    <!-- Don't run immediately at login — let the poller seed fresh data first -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Stdout/stderr for debugging -->
    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/BluePearmain/pipeline.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/BluePearmain/pipeline.log</string>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>__REPO__</string>

    <!-- Don't thrash if it keeps crashing -->
    <key>ThrottleInterval</key>
    <integer>300</integer>
</dict>
</plist>
